# scripts/watchdog/claude-watchdog.service
[Unit]
Description=Claude Code Watchdog v2 - kills hung/leaking Claude processes
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/claude-watchdog.sh
StandardOutput=append:/tmp/claude-watchdog-systemd.log
StandardError=append:/tmp/claude-watchdog-systemd.err
